# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:database-service]
; ── path ชี้ไปที่ /app/database-service/database-service.js ──
command=node /app/database-service/database-service.js
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=1

[program:api-gateway]
; ── path ชี้ไปที่ /app/api-gateway/api-gateway.js ──
command=node /app/api-gateway/api-gateway.js
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
priority=2
; ✅ รอ database-service ขึ้นก่อน 3 วินาที
startsecs=3
